# coredns/coredns.service
[Unit]
Description=CoreDNS

[Service]
Type=notify
ExecStart=/usr/local/bin/coredns -conf /home/nas/coredns/Corefile